# supervisor.conf
[program:my_bbs]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.conf.py
directory=/var/www/my_bbs
autostart=true
autorestart=true